# build/supervisord.all.conf
[supervisord]
nodaemon=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/1
stderr_logfile_maxbytes=0
redirect_stderr=true
loglevel=debug

[program:manager]
command=java -jar /usr/local/bin/admin-assembly-1.0.jar
exitcodes=0,143

[program:monitor]
command=/usr/local/bin/monitor -d
autorestart=false
stopwaitsecs=600
